# home/.config/systemd/user/daily.service
[Unit]
Description=Run daily script

[Service]
Type=oneshot
ExecStart=/home/bruno/.config/systemd/daily
